# Web Stack/deploy/uwsgi.conf
description "ubuntu uwsgi instance"
start on runlevel [2345]
stop on runlevel [06]
exec uwsgi --ini /home/ubuntu/seads/ShR2/Web\ Stack/deploy/uwsgi.ini
